# fixbrightness.conf
description "Sets brightness after graphics device is loaded"

start on graphics-device-added
task
exec /bin/echo 150 > /sys/class/backlight/intel_backlight/brightness
